# src/naga.service
[Unit]
Description=Custom key bindings for the side panel of razer mouses.
After=x_server_started.service

[Service]
ExecStart=/usr/local/bin/naga serviceHelper
User=razerInput

[Install]
WantedBy=default.target